# src/CharmCreator/Templates/Services/Cron.cron
---
name: Crontab entry to run cron:run
fields:
  USERNAME:
    name: Name of PHP / web user who runs the cron job (e.g. www-data)
    type: input
  PATH:
    name: Absolute path to the project dir (e.g. /var/www/myproject)
    type: input
---
# Crontab entry for charm project
* * * * *  USERNAME  php PATH/bob.php cron:run >> /dev/null 2>&1